# deploy/systemd/atlas.service
[Unit]
# Atlas — personal AI Operating System (Sprint 9, ADR-0053).
# Runs the API server, which on startup boots the full kernel: scheduler workers,
# health monitor, ingestion, memory-prune, and scheduled backups. One service is
# enough for a single self-hosted node.
Description=Atlas AI Operating System
After=network-online.target postgresql.service
Wants=network-online.target
# Ollama runs on the host; start it before Atlas if managed by systemd:
# After=network-online.target postgresql.service ollama.service

[Service]
Type=simple
User=atlas
Group=atlas
WorkingDirectory=/data/atlas

# Secrets + host-specific overrides (ATLAS_DB_PASSWORD, ATLAS_API_KEYS, ...).
EnvironmentFile=/etc/atlas/atlas.env

# uv resolves the locked environment; adjust the path to your uv install.
ExecStart=/usr/bin/env uv run atlas serve
Restart=on-failure
RestartSec=5
TimeoutStopSec=30

# --- hardening ---
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
# Atlas writes runtime data + backups here; grant write access.
ReadWritePaths=/data/atlas_data
ProtectKernelTunables=true
ProtectControlGroups=true
RestrictSUIDSGID=true
LockPersonality=true

[Install]
WantedBy=multi-user.target
